# scripts/seagull
#!/sbin/openrc-run

name="Seagull Graphics"
description="Seagull Graphics"

command_user="root"

depends() {
    needs localmount
}

start() {
    xinit seagull -- :0
}
