# scripts/seagull
#!/sbin/openrc-run

name="Seagull Graphics"
description="Seagull Graphics"

command=xinit seagull -- :0

depends() {
    needs localmount
}
